# lact.service
[Unit]
Description=AMDGPU Control Daemon

[Service]
ExecStart=/usr/local/bin/lact-daemon

[Install]
WantedBy=multi-user.target