# res/lactd.service
[Unit]
Description=AMDGPU Control Daemon
After=multi-user.target

[Service]
ExecStart=lact daemon

[Install]
WantedBy=multi-user.target
