# res/lactd.service
[Unit]
Description=AMDGPU Control Daemon
After=multi-user.target

[Service]
ExecStart=lact daemon
Nice=-10

[Install]
WantedBy=multi-user.target
