# deb/usr/lib/systemd/system/lactd.service
[Unit]
Description=AMDGPU Control Daemon
After=multi-user.target

[Service]
ExecStart=/usr/bin/lact-daemon

[Install]
WantedBy=graphical.target
